# docs/HARDWARE/Digital/PEX-1202L/ixpio/ixpio/ixpio-inst.service
[Unit]
Description=IxPIO autoinstall for systemd.

[Service]
Type=oneshot
WorkingDirectory=/root/winson/ixpio
ExecStart=/root/winson/ixpio/ixpio.inst

[Install]
WantedBy=multi-user.target

